# deploy/cointrading-grid-engine.timer
[Unit]
Description=BTCUSDC maker grid lifecycle every minute

[Timer]
OnBootSec=120
OnUnitInactiveSec=60
AccuracySec=10
Unit=cointrading-grid-engine.service

[Install]
WantedBy=timers.target
