# deploy/cointrading-live-supervisor-notify.service
[Unit]
Description=Send live supervisor candidate Telegram alerts

[Service]
Type=oneshot
User=msm1357953_gmail_com
WorkingDirectory=/home/msm1357953_gmail_com/cointrading
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/msm1357953_gmail_com/cointrading/.venv/bin/python -m cointrading.cli live-supervisor-notify
