# deploy/cointrading-strategy-notify.service
[Unit]
Description=Send strategy evaluation Telegram reports

[Service]
Type=oneshot
User=msm1357953_gmail_com
WorkingDirectory=/home/msm1357953_gmail_com/cointrading
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/msm1357953_gmail_com/cointrading/.venv/bin/python -m cointrading.cli strategy-notify
