# deploy/cointrading-trade-event-notify.timer
[Unit]
Description=Check trade events and position status every minute

[Timer]
OnBootSec=300
OnUnitInactiveSec=60
AccuracySec=15
Unit=cointrading-trade-event-notify.service

[Install]
WantedBy=timers.target
