# deploy/cointrading-consecutive-bar-alert.service
[Unit]
Description=Telegram alert when N consecutive same-direction bars print

[Service]
Type=oneshot
User=msm1357953_gmail_com
WorkingDirectory=/home/msm1357953_gmail_com/cointrading
Environment=PYTHONUNBUFFERED=1
TimeoutStartSec=55
ExecStart=/home/msm1357953_gmail_com/cointrading/.venv/bin/python -m cointrading.cli consecutive-bar-alert --symbols BTCUSDC --interval 15m --thresholds 6 7
