# deploy/cointrading-vibe-probe-notify.service
[Unit]
Description=Run Vibe-style research probe and Telegram summary

[Service]
Type=oneshot
User=msm1357953_gmail_com
WorkingDirectory=/home/msm1357953_gmail_com/cointrading
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/msm1357953_gmail_com/cointrading/.venv/bin/python -m cointrading.cli vibe-probe-notify --periodic-minutes 360
